# omx-bridge.service
[Unit]
Description=OMX Bridge Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/omx-bridge
ExecStart=/usr/bin/env node dist/main.js
Restart=always
RestartSec=5
EnvironmentFile=%h/workspace/omx-bridge/.env
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
